# main/systemd/longtail.service
[Unit]
Description=Long-Tail Arbitrage Engine
After=network.target

[Service]
Type=simple
User=longtail
WorkingDirectory=/opt/longtail
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/longtail/venv/bin/python -m scripts.run_live
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=longtail

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/longtail

[Install]
WantedBy=multi-user.target
